# app-misc/bottle-time-processor/files/bottle-time-processor.service
[Unit]
Description=Bottle Time Processor Daemon
After=network.target

[Service]
# The daemon runs in the foreground, so use Type=simple.
Type=simple
# Adjust the path to the binary as needed.
ExecStart=/usr/bin/bottle-time-processor

# Enable a watchdog to monitor the service.
WatchdogSec=900

# Restart the service automatically on failure.
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target